# systemd/repdojo.target
[Unit]
Description=RepDojo Service Group
Requires=repdojo-web.service repdojo-api.service
After=network.target

[Install]
WantedBy=multi-user.target
